# corefreqd.service
[Unit]
Description=CoreFreq Daemon
After=systemd-modules-load.service dkms.service
Wants=dkms.service

[Service]
Type=simple
# 简洁明了 - 如果模块未就绪，服务失败，systemd 会自动重试
ExecStart=/usr/bin/corefreqd -q
ExecStopPost=-/sbin/modprobe -r corefreqk
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target