# corefreqd.service
[Unit]
Description=CoreFreq Daemon
After=systemd-modules-load.service dkms.service
Wants=dkms.service

[Service]
Type=simple
# 确保模块加载后再启动守护进程
ExecStartPre=/sbin/modprobe corefreqk
ExecStart=/usr/bin/corefreqd -q
ExecStopPost=-/sbin/modprobe -r corefreqk
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target